# rpm/SOURCES/vexor-victorialogs.service
[Unit]
Description=Vexor / VictoriaLogs daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vexor
Group=vexor
EnvironmentFile=-/etc/vexor/logs.env
ExecStart=/usr/bin/vexor-victorialogs
Restart=on-failure
RestartSec=5s
LimitNOFILE=1048576
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
